# src/testbed-tools/templates/iperf.supervisor.conf
[program:iperfserver]
command=/usr/bin/iperf -s
autostart=true
autorestart=true
stderr_logfile=/mnt/nfs/results/iperf.err.log
stdout_logfile=/mnt/nfs/results/iperf.out.log
